# debian/cron.d
#
# Regular cron jobs for the bbbs-equivs package
#
0 4	* * *	root	bbbs-equivs_maintenance
