# debian/cron.d
#
# Regular cron jobs for the bbbs-installer package
#
0 4	* * *	root	bbbs-installer_maintenance
